# templates/button_handler.service.j2
[Unit]
Description=button_handler
After=network.target

[Service]
Type=simple
ExecStart=/var/lib/doorbell/button_handler
PIDFile=/var/run/button_handler.pid

[Install]
WantedBy=multi-user.target
